# build/build/cfg/etc/sv/powerpanel/run
#!/bin/sh
exec 2>&1
exec chpst -u root /usr/local/bin/pwrstatd 1>&2
