# build/build/cfg/etc/udev/rules.d/60-cdrom.rules
KERNEL=="sr0", SYMLINK+="cdrom", RUN+="chmod -R g+r /dev/sr0"
